# scripts/crontab
@reboot /bin/clean.sh >> /var/log/pacman-clean.log
* 5 * * * /bin/clean.sh >> /var/log/pacman-clean.log
