# scripts/crontab
@reboot /bin/update-and-clean.sh >> /var/log/pacman-update-clean.log
* 5 * * * /bin/update-and-clean.sh >> /var/log/pacman-update-clean.log
